# pi/pbai-heatdrip.service
[Unit]
Description=PBAI Heat Drip (keeps Ego alive for sustained play)
After=pbai-daemon.service
Requires=pbai-daemon.service

[Service]
Type=simple
User=pbai
ExecStart=/bin/bash -c 'while true; do curl -s -X POST http://localhost:8420/inject/heat -H "Content-Type: application/json" -d "{\"target\":\"ego\",\"amount\":10}" > /dev/null; curl -s -X POST http://localhost:8420/inject/heat -H "Content-Type: application/json" -d "{\"target\":\"identity\",\"amount\":6}" > /dev/null; sleep 3; done'
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
